# supervisor-uwsgi.conf
[program:uwsgi]
command = uwsgi --ini /usr/src/api/uwsgi.ini
stopsignal=QUIT  
stopasgroup=true 
killasgroup=true 
autostart=true   
autorestart=true
redirect_stderr=true 
logfile_backups=5